# root_conf/etc/NetworkManager/dispatcher.d/98-thetoto.sh
/etc/NetworkManager/dispatcher.d/98-thetoto.sh